# roles/manager/files/bluechi-controller.socket
[Socket]
ListenStream=/var/lib/haproxy/bluechi-controller.sock

[Install]
WantedBy=multi-user.target
