# roles/manager/files/hirte.socket
[Socket]
ListenStream=/var/lib/haproxy/bluechi-manager.sock

[Install]
WantedBy=multi-user.target
