# roles/manager/files/hirte.socket
[Socket]
ListenStream=/var/lib/haproxy/hirte.sock

[Install]
WantedBy=multi-user.target
